# electricity-monitor.service
[Unit]
Description=Electricity Monitor Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=debian
Group=debian
WorkingDirectory=/opt/electricity-monitor
ExecStart=/opt/electricity-monitor/electricity-monitor

EnvironmentFile=/etc/electricity-monitor.env

# Restart policy
Restart=always
# RestartSec=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/electricity-monitor /var/lib/electricity-monitor

# Resource limits
LimitNOFILE=1024

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=electricity-monitor

# Required for ICMP ping
AmbientCapabilities=CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_RAW

[Install]
WantedBy=multi-user.target
